# sample_biocs_fixed.json
[
  {
    "name": "Suspicious PowerShell Execution",
    "comment": "Encoded PowerShell execution detection",
    "description": "Detect suspicious PowerShell encoded commands",
    "severity": "SEV_040_HIGH",
    "status": "ENABLED",
    "is_xql": true,
    "category": "EXECUTION",
    "xql": "dataset = xdr_data | filter action_process_image_name contains \"powershell.exe\" | filter action_process_image_command_line contains \"-enc\" | limit 100",
    "indicator_text": "dataset = xdr_data | filter action_process_image_name contains \"powershell.exe\" | filter action_process_image_command_line contains \"-enc\" | limit 100",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {}
    },
    "mitre_technique_id_and_name": "T1059.001 - PowerShell",
    "mitre_technique_id": "T1059.001",
    "mitre_tactic_id_and_name": "TA0002 - Execution",
    "mitre_tactic_id": "TA0002"
  },
  {
    "name": "Administrator File Execution",
    "comment": "Admin file execution detection",
    "description": "Detect administrator account executing files",
    "severity": "SEV_030_MEDIUM",
    "status": "ENABLED",
    "is_xql": true,
    "category": "PRIVILEGE_ESCALATION",
    "xql": "dataset = xdr_data | filter action_process_username =~ \"Administrator\" | filter event_type = PROCESS_START | limit 100",
    "indicator_text": "dataset = xdr_data | filter action_process_username =~ \"Administrator\" | filter event_type = PROCESS_START | limit 100",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {}
    },
    "mitre_technique_id_and_name": "T1078.001 - Valid Accounts Default Accounts",
    "mitre_technique_id": "T1078.001",
    "mitre_tactic_id_and_name": "TA0004 - Privilege Escalation",
    "mitre_tactic_id": "TA0004"
  },
  {
    "name": "Malicious Domain Connection",
    "comment": "malware.com",
    "description": "Detect connections to known malicious domains",
    "severity": "SEV_050_CRITICAL",
    "status": "ENABLED",
    "is_xql": false,
    "category": "TA0011",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {
        "PROCESS_EXECUTION_EVENT": {
          "filter": {
            "AND": [
              {
                "SEARCH_FIELD": "network_event",
                "SEARCH_TYPE": "CONTAINS",
                "SEARCH_VALUE": "contains",
                "EXTRA_FIELDS": [],
                "isExtended": false,
                "node": "attributes"
              }
            ]
          }
        }
      }
    },
    "indicator_text": "network_event contains contains",
    "mitre_technique_id_and_name": "Known malware domain communication - Unknown",
    "mitre_technique_id": "Known malware domain communication",
    "mitre_tactic_id_and_name": "T1071 - Application Layer Protocol",
    "mitre_tactic_id": "T1071"
  },
  {
    "name": "Suspicious Process Execution",
    "comment": "C:\\Users\\Public\\",
    "description": "Detect execution from suspicious directories",
    "severity": "SEV_040_HIGH",
    "status": "ENABLED",
    "is_xql": false,
    "category": "TA0002",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {
        "PROCESS_EXECUTION_EVENT": {
          "filter": {
            "AND": [
              {
                "SEARCH_FIELD": "process_execution_event",
                "SEARCH_TYPE": "CONTAINS",
                "SEARCH_VALUE": "contains",
                "EXTRA_FIELDS": [],
                "isExtended": false,
                "node": "attributes"
              }
            ]
          }
        }
      }
    },
    "indicator_text": "process_execution_event contains contains",
    "mitre_technique_id_and_name": "Execution from public directory - Unknown",
    "mitre_technique_id": "Execution from public directory",
    "mitre_tactic_id_and_name": "T1059 - Command and Scripting Interpreter",
    "mitre_tactic_id": "T1059"
  },
  {
    "name": "Registry Persistence",
    "comment": "SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run",
    "description": "Monitor registry modifications for persistence",
    "severity": "SEV_030_MEDIUM",
    "status": "ENABLED",
    "is_xql": false,
    "category": "TA0003",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {
        "PROCESS_EXECUTION_EVENT": {
          "filter": {
            "AND": [
              {
                "SEARCH_FIELD": "registry_event",
                "SEARCH_TYPE": "CONTAINS",
                "SEARCH_VALUE": "contains",
                "EXTRA_FIELDS": [],
                "isExtended": false,
                "node": "attributes"
              }
            ]
          }
        }
      }
    },
    "indicator_text": "registry_event contains contains",
    "mitre_technique_id_and_name": "Registry run key persistence - Unknown",
    "mitre_technique_id": "Registry run key persistence",
    "mitre_tactic_id_and_name": "T1547 - Boot or Logon Autostart Execution",
    "mitre_tactic_id": "T1547"
  },
  {
    "name": "Fileless Attack Detection",
    "comment": "Fileless PowerShell execution",
    "description": "Detect fileless script loading activity",
    "severity": "SEV_040_HIGH",
    "status": "ENABLED",
    "is_xql": true,
    "category": "EXECUTION",
    "xql": "dataset = xdr_data | filter event_type = MODULE_LOAD | filter action_file_name contains \"powershell.exe\" | filter action_process_image_command_line contains \"-nop -w hidden\" | limit 100",
    "indicator_text": "dataset = xdr_data | filter event_type = MODULE_LOAD | filter action_file_name contains \"powershell.exe\" | filter action_process_image_command_line contains \"-nop -w hidden\" | limit 100",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {}
    },
    "mitre_technique_id_and_name": "T1059.001 - PowerShell",
    "mitre_technique_id": "T1059.001",
    "mitre_tactic_id_and_name": "TA0002 - Execution",
    "mitre_tactic_id": "TA0002"
  },
  {
    "name": "Suspicious Download",
    "comment": "C:\\Users\\*\\Downloads\\",
    "description": "Detect files downloaded from external sources",
    "severity": "SEV_020_LOW",
    "status": "ENABLED",
    "is_xql": false,
    "category": "TA0011",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {
        "PROCESS_EXECUTION_EVENT": {
          "filter": {
            "AND": [
              {
                "SEARCH_FIELD": "file_event",
                "SEARCH_TYPE": "CONTAINS",
                "SEARCH_VALUE": "contains",
                "EXTRA_FIELDS": [],
                "isExtended": false,
                "node": "attributes"
              }
            ]
          }
        }
      }
    },
    "indicator_text": "file_event contains contains",
    "mitre_technique_id_and_name": "File download detection - Unknown",
    "mitre_technique_id": "File download detection",
    "mitre_tactic_id_and_name": "T1105 - Ingress Tool Transfer",
    "mitre_tactic_id": "T1105"
  },
  {
    "name": "Network Beacon Detection",
    "comment": ".*\\.dyn.*\\..*",
    "description": "Detect periodic network beaconing behavior",
    "severity": "SEV_040_HIGH",
    "status": "ENABLED",
    "is_xql": false,
    "category": "TA0011",
    "indicator": {
      "runOnCGO": true,
      "investigationType": "PROCESS_EXECUTION_EVENT",
      "investigation": {
        "PROCESS_EXECUTION_EVENT": {
          "filter": {
            "AND": [
              {
                "SEARCH_FIELD": "network_event",
                "SEARCH_TYPE": "CONTAINS",
                "SEARCH_VALUE": "regex",
                "EXTRA_FIELDS": [],
                "isExtended": false,
                "node": "attributes"
              }
            ]
          }
        }
      }
    },
    "indicator_text": "network_event contains regex",
    "mitre_technique_id_and_name": "Dynamic DNS beacon detection - Unknown",
    "mitre_technique_id": "Dynamic DNS beacon detection",
    "mitre_tactic_id_and_name": "T1071.004 - Application Layer Protocol: DNS",
    "mitre_tactic_id": "T1071.004"
  }
]